# modules/linux_gui/templates/gnome-session-x11.service.erb
[Unit]
Description=Headless GNOME session on X11
After=graphical.target
Wants=graphical.target

[Service]
Environment=DISPLAY=:0
Environment=XDG_SESSION_TYPE=x11
Environment=XDG_RUNTIME_DIR=/run/user/<%= @builder_uid %>
ExecStartPre=/bin/mkdir -p /run/user/<%= @builder_uid %>
ExecStartPre=/bin/chown <%= @builder_uid %>:<%= @builder_gid %> /run/user/<%= @builder_uid %>
ExecStartPre=/bin/chmod 700 /run/user/<%= @builder_uid %>
ExecStart=/usr/bin/dbus-run-session -- gnome-session --session=gnome
Restart=always

[Install]
WantedBy=default.target
