# modules/linux_papertrail/templates/papertrail-syslog.service.erb
[Unit]
Description=Papertrail-syslog
After=systemd-journald.service
Requires=systemd-journald.service
[Service]
ExecStart=/bin/sh -c "journalctl <% @syslog_identifiers.each do |identifier| %>-t <%= identifier %> <% end %> -f | ncat --ssl <%= @papertrail_host %> <%= @papertrail_port %>"
TimeoutStartSec=0
Restart=on-failure
RestartSec=5s
[Install]
WantedBy=multi-user.target
