# modules/linux_generic_worker/files/check_gw.service
[Unit]
Description=check_gw

[Service]
Type=oneshot
ExecStart=/opt/relops-check_gw/check_gw.py --reboot
